# supervisord.conf
[supervisord]
nodaemon=true

[program:docker]
command=/dind.sh

[program:gocd-agent]
command=/gocd-agent.sh

